# supervisor.conf
[program:htg_backend]
command=/data/htg_backend/bvenv/bin/python /data/htg_backend/run.py
autostart=true
autorestart=true
redirect_stderr=false
stopasgroup=true